# veilbot-admin.service
[Unit]
Description=VeilBot Admin Panel
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8000
WorkingDirectory=/root/veilbot/admin
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 